# yt-extractor.service
[Unit]
Description=YouTube Playlist URL Extractor
After=network.target

[Service]
User=www
Group=www
WorkingDirectory=/www/wwwroot/urls_playlists
ExecStart=/www/wwwroot/urls_playlists/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile /var/log/yt-extractor-access.log \
    --error-logfile /var/log/yt-extractor-error.log \
    app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
